# shadowsocks.ini
[program:shadowsocks]
command=python server.py m
process_name=shadowsocks
directory=/usr/local/shadowsocks
startretries=99
user=root
